# extra/etc/supervisor/conf.d/jabe.conf
[program:jabe]
command=/home/melaman/JABE/venv/bin/gunicorn wsgi:app -c /home/melaman/JABE/gunicorn_conf.py
directory=/home/melaman/JABE
user=melaman
autorestart=true
redirect_stderr=true
